# root/etc/s6-overlay/s6-rc.d/radicale-fixattrs/run
#!/command/with-contenv sh

chown -R appuser:appuser /config /opt/nzbhydra2
chmod -R g+w /config /opt/nzbhydra2
